# chatty-ai-launcher.service
[Service]
Type=simple
User=nickspi5
Group=nickspi5
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/nickspi5/.Xauthority"
Environment="HOME=/home/nickspi5"

# Wait for display to be ready
ExecStartPre=/bin/bash -c 'while ! xset q &>/dev/null; do sleep 2; done'

# Run the startup script
ExecStart=/home/nickspi5/Chatty_AI/chatty_ai_startup.sh

# Restart if it fails
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Run with normal priority
Nice=0

# Security settings
NoNewPrivileges=false
PrivateTmp=false

[Install]
WantedBy=graphical.target